# Python/script_alarm
### BEGIN INIT INFO
# Provides:          scipt_alarm
# Required-Start:  $all 
# Required-Stop:  $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: alarm script
# Description:       Run Raspberry Pi alarm service
### END INIT INFO

cd /home/pi/Documents/projekty
sudo python kontrolka_alarmu.py 1800 &
